# systemd-watcher.service
[Unit]
Description=Systemd-Watcher

[Service]
ExecStart=/usr/local/bin/systemd-watcher

[Install]
WantedBy=multi-user.target
